# etc/supervisor/conf.d/riemann-dash.conf
[program:riemann-dash]
command=riemann-dash /etc/riemann/dash.rb
redirect_stderr=true
