# etc/supervisor/conf.d/tcollector.conf
[program:tcollector]
command=/opt/tcollector/tcollector.py --host 192.168.12.10
