# etc/supervisor/conf.d/zookeeper.conf
[program:zookeeper]
command=/opt/zookeeper/bin/zkServer.sh start-foreground
directory=/opt/zookeeper
user=zookeeper
redirect_stderr=true
