# etc/init/kibana.conf
description "kibana"

start on virtual-filesystems
stop on runlevel [06]

respawn
respawn limit 5 10

console log

setuid nobody

exec /opt/kibana/bin/kibana -c /etc/kibana/kibana.yml
