# etc/init/oculus-webapp.conf
description "oculus-webapp"

start on virtual-filesystems
stop on runlevel [06]

console log

pre-start script
  mkdir -p /var/run/oculus
  chown oculus:oculus /var/run/oculus
end script

exec start-stop-daemon --start -c oculus -d /opt/oculus --exec /usr/local/bin/bundle exec thin start
