# etc/init/skyline-analyzer.conf
description "skyline-analyzer"

start on virtual-filesystems
stop on runlevel [06]

respawn
respawn limit 5 10

console log

chdir /opt/skyline
setuid skyline
setgid skyline

exec bin/analyzer.d run
